# ctl/90-asus-nb-acpi.rules
ACTION=="change", SUBSYSTEM=="power_supply", ATTR{online}=="*", RUN+="/usr/bin/systemctl restart asus-nb-acpi-d.service"
